# v2-interface/packages/api/scripts/deploy/kasane-v2-user-positions/kasane-v2-user-positions-indexer.service
[Unit]
Description=Kasane Indexer (token/pool + optional user LP)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/kasane-dex/v2-interface/packages/api
EnvironmentFile=/etc/kasane-indexer.env
ExecStart=/usr/bin/env bun run token-indexer:run
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
